# docker/lsio/root/etc/services.d/harmony/run
#!/usr/bin/with-contenv bash
set -euo pipefail

PUID="${PUID:-911}"
PGID="${PGID:-911}"
TZ="${TZ:-Etc/UTC}"
UMASK_VALUE="${UMASK:-022}"

export TZ

# Ensure abc group matches requested PGID
if getent passwd abc >/dev/null 2>&1; then
    CURRENT_UID="$(id -u abc)"
    if [ "$CURRENT_UID" != "$PUID" ]; then
        deluser abc >/dev/null 2>&1 || true
    fi
fi

if getent group abc >/dev/null 2>&1; then
    CURRENT_GID="$(getent group abc | cut -d: -f3)"
    if [ "$CURRENT_GID" != "$PGID" ]; then
        delgroup abc >/dev/null 2>&1 || true
    fi
fi

if ! getent group abc >/dev/null 2>&1; then
    addgroup -S -g "$PGID" abc
fi

if ! getent passwd abc >/dev/null 2>&1; then
    adduser -S -D -h /config -s /sbin/nologin -u "$PUID" -G abc abc
fi

mkdir -p /config /downloads /music
chown -R abc:abc /config /downloads /music

# Apply requested umask
if [ -n "$UMASK_VALUE" ]; then
    umask "$UMASK_VALUE"
fi

export SQLALCHEMY_DATABASE_URL="${SQLALCHEMY_DATABASE_URL:-sqlite:////config/harmony.db}"
export HOME=/config

if command -v alembic >/dev/null 2>&1 && [ -f /app/alembic.ini ]; then
    echo "[harmony] Running Alembic migrations"
    s6-setuidgid abc alembic upgrade head
fi

echo "[harmony] Starting uvicorn on 0.0.0.0:8080"
exec s6-setuidgid abc python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080
