# docker/services.d/harmony/run
#!/command/with-contenv sh
set -eu

cd /app/app
exec /app/.venv/bin/python -m app.runtime.container_entrypoint "$@"
