# docker/services.d/harmony/run
#!/command/with-contenv bash
set -euo pipefail
export PATH="/app/.venv/bin:${PATH}"
exec python -m app.runtime.container_entrypoint
